# crontab
# Clean ~/tmp
0 * * * * find ~/tmp -mtime +7 -delete
